# tarsnap.service
[Unit]
Description=Tarsnap service
Documentation=man:tarsnap

[Service]
ExecStart=/home/loy/tools/tarsnap-backup.sh
Restart=always
