# ip_mailer/ipmailer.timer
[Unit]
Description=Check external IP address daily

[Timer]
OnCalendar=hour
Persistent=true     
 
[Install]
WantedBy=timers.target
